# systemd/vibhost-agent.service
[Unit]
Description=VibHost Resident gRPC Agent
Documentation=https://github.com/Aventus-Group/vh-agent
After=network-online.target wg-quick@vhnet0.service
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/vh-agent
Environment=VH_AGENT_CONFIG=/etc/vibhost/agent.conf
Restart=always
RestartSec=3
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home /tmp /var/log /etc/nginx /etc/ssl /etc/vibhost
PrivateTmp=false
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
