# systemd/vibhost-agent.service
[Unit]
Description=VibHost Host Agent
After=network-online.target wg-quick@vhnet0.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vh-agent
Restart=always
RestartSec=10
User=root
Group=root

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/usr/local/bin /root/.vibhost/wg /etc/vibhost /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
